# com.pishleme.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pishleme.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/pishleme</string>
        <string>--app</string>
        <string>Safari</string>
        <string>--time</string>
        <string>3600</string>
        <string>--app</string>
        <string>Discord</string>
        <string>--time</string>
        <string>1800</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/pishleme.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/pishleme.error.log</string>
    <key>UserName</key>
    <string>root</string>
    <key>GroupName</key>
    <string>wheel</string>
    <key>ProcessType</key>
    <string>Background</string>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>